# .config/systemd/user/flatpak-update.service
[Unit]
Description=Update flatpaks

[Service]
ExecStart=flatpak update --noninteractive
