# .config/systemd/user/sort-traderepublic.service
[Unit]
Description=Run sort-traderepublic

[Service]
Type=oneshot
ExecStart=/usr/local/bin/sort-traderepublic /home/jak/Documents/Bank/TradeRepublic/
