# .config/systemd/user/sort-traderepublic.path
[Unit] 
Wants=sort-traderepublic.service 

[Path] 
PathChanged=/home/jak/Documents/Bank/TradeRepublic/
